# delayBoot1.service
#!/bin/sh

[Unit]
Description=Run ARDPI
After=network.target
 
[Service]
ExecStart=/bin/bash /home/pi/apl/runBATMAN.sh
 
[Install]
WantedBy=multi-user.target



