# delayBoot3.service
#!/bin/sh

[Unit]
Description=Run ARDPI
After=network.target
 
[Service]
ExecStart=/bin/bash /home/pi/apl3/runBATMAN3.sh
 
[Install]
WantedBy=multi-user.target
